# udev-firmware.rules
ACTION=="add", SUBSYSTEM=="firmware", RUN+="/sbin/firmware_helper"
